# luci-app-campusnet/root/etc/uci-defaults/luci-campusnet
uci -q batch <<-EOF >/dev/null
[ ! -f "/usr/share/ucitrack/luci-app-campusnet.json" ] && {
    cat > /usr/share/ucitrack/luci-app-campusnet.json << EEOF
{
    "config": "campusnet",
    "init": "campusnet"
}
EEOF
}
EOF

rm -f /tmp/luci-indexcache
exit 0
