# luci-app-nezhav1/root/etc/uci-defaults/luci-nezha-agent
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-nezha-agent.json" ] && {
    cat > /usr/share/ucitrack/luci-app-nezha-agent.json << EEOF
{
    "config": "nezha-agent",
    "init": "nezha-agent"
}
EEOF
}

# 设置权限
chmod +x /etc/init.d/nezha-agent

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@nezha-agent[-1]
	add ucitrack nezha-agent
	set ucitrack.@nezha-agent[-1].init=nezha-agent
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
