# luci-app-thunder/root/etc/uci-defaults/luci-thunder
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@thunder[-1]
	add ucitrack thunder
	set ucitrack.@thunder[-1].init=thunder
	commit ucitrack
EOF

rm -rf /tmp/luci-*
exit 0
