# luci-app-oscam/root/etc/uci-defaults/luci-oscam
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-oscam.json" ] && {
    cat > /usr/share/ucitrack/luci-app-oscam.json << EEOF
{
    "config": "oscam",
    "init": "oscam"
}
EEOF
}

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@oscam[-1]
	add ucitrack oscam
	set ucitrack.@oscam[-1].init=oscam
	commit ucitrack
EOF


/etc/init.d/oscam stop

rm -f /tmp/luci-indexcache
exit 0
