# luci-app-aliyundrive-webdav/root/etc/uci-defaults/luci-aliyundrive-webdav
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-aliyundrive-webdav.json" ] && {
    cat > /usr/share/ucitrack/luci-app-aliyundrive-webdav.json << EEOF
{
    "config": "aliyundrive-webdav",
    "init": "aliyundrive-webdav"
}
EEOF
}
/etc/init.d/ucitrack reload

uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
