# luci-app-adguardhome/root/etc/uci-defaults/luci-adguardhome
#!/bin/sh

chmod +x /etc/init.d/adguardhome /usr/libexec/rpcd/luci.adguardhome

exit 0
