# luci-app-caddy/root/etc/uci-defaults/luci-caddy
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-caddy.json" ] && {
    cat > /usr/share/ucitrack/luci-app-caddy.json << EEOF
{
    "config": "caddy",
    "init": "caddy"
}
EEOF
}
/etc/init.d/ucitrack reload

uci -q batch <<-EOF >/dev/null
EOF

rm -rf /tmp/luci-*
chmod 755 /etc/init.d/caddy
exit 0
