# luci-app-control-timewol/root/etc/uci-defaults/luci-app-control-timewol
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-timewol.json" ] && {
    cat > /usr/share/ucitrack/luci-app-timewol.json << EEOF
{
    "config": "timewol",
    "init": "timewol"
}
EEOF
}

uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
