# luci-app-shadowsocks-libev/root/etc/uci-defaults/40_luci-shadowsocks-libev
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-shadowsocks-libev.json" ] && {
    cat > /usr/share/ucitrack/luci-app-shadowsocks-libev.json << EEOF
{
    "config": "shadowsocks-libev",
    "init": "shadowsocks-libev"
}
EEOF
/etc/init.d/ucitrack reload
}

uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
