# luci-app-switchtools/root/etc/uci-defaults/40_luci-switchtools
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-switchmode.json" ] && {
    cat > /usr/share/ucitrack/luci-app-switchmode.json << EEOF
{
    "config": "switchmode",
    "init": "switchmode"
}
EEOF
}

uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
