# .github/diy/packages/luci-app-pptpd/root/etc/uci-defaults/luci-pptpd
#!/bin/sh

[ ! -f "/usr/share/ucitrack/luci-app-pptpd.json" ] && {
    cat > /usr/share/ucitrack/luci-app-pptpd.json << EEOF
{
    "config": "pptpd",
    "init": "pptpd"
}
EEOF
}

chmod +x /etc/init.d/pptpd
rm -f /tmp/luci-indexcahe
exit 0
