# luci-app-netkeeper-interception/root/etc/uci-defaults/40_luci-app-netkeeper-interception
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-netkeeper-interception.json" ] && {
    cat > /usr/share/ucitrack/luci-app-netkeeper-interception.json << EEOF
{
    "config": "netkeeper-interception",
    "init": "netkeeper-interception"
}
EEOF
/etc/init.d/ucitrack reload
}
uci -q batch <<-EOF >/dev/null
EOF
rm -rf /var/luci-modulecache/; rm -f /var/luci-indexcache;
exit 0

