# luci-app-log/root/uci-defaults/luci-log
#!/bin/sh

sed -i "s/luci-mod-status-logs/luci-mod-status-log/" /usr/share/luci/menu.d/luci-mod-status.json

exit 0
