# luci-app-dnsmasq-ipset/files/root/etc/uci-defaults/luci-dnsmasq-ipset
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-dnsmasq-ipset.json" ] && {
    cat > /usr/share/ucitrack/luci-app-dnsmasq-ipset.json << EEOF
{
    "config": "dnsmasq-ipset",
    "init": "dnsmasq-ipset"
}
EEOF
/etc/init.d/ucitrack reload
}
uci -q batch <<-EOF >/dev/null
EOF
exit 0
