# my-default-settings/target/linux/mediatek/filogic/base-files/etc/hotplug.d/iface/99-tr3000
#!/bin/sh

board=$(board_name)

[ "$board" = "cudy,tr3000-v1" ] && \
[ "$ACTION" = "ifup" ] && \
ethtool -s wan autoneg on
