# luci-app-subconverter/root/etc/uci-defaults/luci-subconverter
#!/bin/sh
chmod +x /etc/init.d/subconverter
exit 0