# luci-app-npc/root/etc/uci-defaults/luci-npc
#!/bin/sh
[ ! -f "/usr/share/ucitrack/luci-app-npc.json" ] && {
    cat > /usr/share/ucitrack/luci-app-npc.json << EEOF
{
    "config": "npc",
    "init": "npc"
}
EEOF
}

uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
