# luci-app-subconverter/root/etc/init.d/subconverter
#!/bin/sh /etc/rc.common

START=90
STOP=10

start() {
  /usr/subconverter/subconverter > /dev/null 2>&1 &
}

stop() {
  killall -q subconverter
}

restart() {
  stop
  start
}
