# luci-app-macvlan/files/root/etc/uci-defaults/40_luci-app-macvlan

[ ! -f "/usr/share/ucitrack/luci-app-macvlan.json" ] && {
    cat > /usr/share/ucitrack/luci-app-macvlan.json << EEOF
{
    "config": "macvlan",
    "init": "macvlan"
}
EEOF
}
/etc/init.d/ucitrack reload
uci -q batch <<-EOF >/dev/null
EOF

rm -f /tmp/luci-indexcache
exit 0
